# deploy/helpdesk-backup.timer
[Unit]
Description=Daily Corporate Help Desk backup

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
