# deploy/helpdesk.service
[Unit]
Description=Corporate Help Desk Application
After=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/helpdesk
EnvironmentFile=/home/pi/helpdesk/.env
ExecStart=/home/pi/helpdesk/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
